# deploy/naozhi.service
[Unit]
Description=Naozhi - Claude Code IM Gateway
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify + WatchdogSec: naozhi (main.go) calls sd_notify("READY=1")
# when the listener is bound and pings periodically for liveness. Using
# Type=simple here produced an infinite "sd_notify READY failed" log
# loop and disabled watchdog-driven restart on hang.
Type=notify
NotifyAccess=main
WatchdogSec=120
# TODO: Change User, WorkingDirectory, and ExecStart paths to match your deployment
User=ec2-user
WorkingDirectory=/home/ec2-user/naozhi
ExecStart=/home/ec2-user/naozhi/bin/naozhi --config /home/ec2-user/naozhi/config.yaml
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=5

# Shim process survival: only kill main process on stop, never SIGKILL the cgroup.
# Shims are moved to /sys/fs/cgroup/naozhi-shims/ and must survive service restarts.
KillMode=process
SendSIGKILL=no
TimeoutStopSec=5

# Runtime environment
Environment=HOME=/home/ec2-user
Environment=PATH=/home/ec2-user/.local/bin:/home/ec2-user/.nvm/versions/node/v22.22.1/bin:/usr/local/bin:/usr/bin:/bin

# Bedrock/Anthropic env is read from ~/.claude/settings.json by Claude CLI.
# Do not set ANTHROPIC_* / AWS_REGION / CLAUDE_CODE_USE_BEDROCK here — keep a single
# source of truth in settings.json to match TUI behavior.

# Feishu credentials (create with deploy/setup-env.sh)
EnvironmentFile=-/home/ec2-user/.naozhi/env

[Install]
WantedBy=multi-user.target
