# deploy/naozhi.service
# RNEW-OPS-412: This file is a manual-deploy reference that stays in sync
# with cmd/naozhi/service.go::generateSystemdUnit. The Go generator is the
# authoritative source used by `naozhi install`; this template exists for
# operators who prefer `cp` + edit in place without running the binary.
# If you change this file, update generateSystemdUnit and its
# TestGenerateSystemdUnit_MatchesDeployTemplate regression test to match.
# Prefer `sudo ./naozhi install` for production deploys.

[Unit]
Description=Naozhi - Claude Code IM Gateway
After=network-online.target
Wants=network-online.target
# StartLimit* belong in [Unit] (systemd v230+); in [Service] they are ignored
# with a warning, so crash-loop rate-limiting silently never triggers.
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
# Type=notify + WatchdogSec: naozhi (main.go) calls sd_notify("READY=1")
# when the listener is bound and pings periodically for liveness. Using
# Type=simple here produced an infinite "sd_notify READY failed" log
# loop and disabled watchdog-driven restart on hang.
Type=notify
NotifyAccess=main
WatchdogSec=120
# TODO: Change User, WorkingDirectory, and ExecStart paths to match your deployment
User=ec2-user
WorkingDirectory=/home/ec2-user/naozhi
ExecStart=/home/ec2-user/naozhi/bin/naozhi --config /home/ec2-user/naozhi/config.yaml
Restart=always
RestartSec=5

# Shim process survival: only kill main process on stop, never SIGKILL the cgroup.
# Shims are moved to /sys/fs/cgroup/naozhi-shims/ and must survive service restarts.
KillMode=process
SendSIGKILL=no
TimeoutStopSec=5

# Runtime environment
Environment=HOME=/home/ec2-user
Environment=PATH=/home/ec2-user/.local/bin:/home/ec2-user/.nvm/versions/node/v22.22.1/bin:/usr/local/bin:/usr/bin:/bin

# Bedrock/Anthropic env is read from ~/.claude/settings.json by Claude CLI.
# Do not set ANTHROPIC_* / AWS_REGION / CLAUDE_CODE_USE_BEDROCK here — keep a single
# source of truth in settings.json to match TUI behavior.

# Feishu credentials (create with deploy/setup-env.sh)
EnvironmentFile=-/home/ec2-user/.naozhi/env

[Install]
WantedBy=multi-user.target
